# bin/x86_64/usr/share/ctdb/tests/eventscripts/etc-ctdb/rc.local
version https://git-lfs.github.com/spec/v1
oid sha256:0526cc2cccbfa0a7a9c5aebac6ad1565d6ef98a5d47a5ad43740044116026a3c
size 1140
