# bin/x86_64/etc/inittab
version https://git-lfs.github.com/spec/v1
oid sha256:1291f880b9de3e1774670e85aa282575c70b1d0c26ade66980c59b516ddf2f58
size 1238
